# config/homebrain-health.timer
[Unit]
Description=Run HomeBrain health check every 30 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
RandomizedDelaySec=2min

[Install]
WantedBy=timers.target
