# config/llama-server.service
[Unit]
Description=llama.cpp inference server
After=network.target multi-user.target
Wants=multi-user.target

[Service]
Type=simple
User=__HOMEBRAIN_USER__
Group=__HOMEBRAIN_USER__
WorkingDirectory=__WORKDIR__
Environment="LD_LIBRARY_PATH=__WORKDIR__"
Environment="GGML_VK_DEVICE=0"

# Wait for GPU device node before starting inference
ExecStartPre=/bin/sh -c 'until [ -e /dev/dri/card0 ] || [ -e /dev/dri/renderD128 ]; do sleep 1; done'

ExecStart=__LLAMA_BIN__ \
  --model __LLAMA_MODEL__ \
  --ctx-size __CTX_SIZE__ \
  --host 127.0.0.1 \
  --port 8001 \
  -ngl __NGL__ __EXTRA_FLAGS__

Restart=always
RestartSec=10
TimeoutStartSec=600
OOMScoreAdjust=-500
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
