# config/llama-server.service
[Unit]
Description=llama.cpp inference server
After=network.target multi-user.target
Wants=multi-user.target

[Service]
Type=simple
User=admin
Group=admin
WorkingDirectory=@BIN_DIR@
Environment="LD_LIBRARY_PATH=@BIN_DIR@"
Environment="GGML_VK_DEVICE=0"

# Wait for GPU device to be available before starting
ExecStartPre=/bin/sh -c 'until [ -e /dev/dri/card0 ] || [ -e /dev/dri/renderD128 ]; do sleep 1; done'

ExecStart=@BIN_PATH@ --model @MODEL_PATH@ --ctx-size @CTX_SIZE@ --host 127.0.0.1 --port 8001 -ngl @NGL@@EXTRA_FLAGS@

Restart=always
RestartSec=10
TimeoutStartSec=600
OOMScoreAdjust=-500
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
